# benchmark-service/deploy/systemd/qwen-benchmark-worker.service
[Unit]
Description=Qwen Code benchmark worker
After=network-online.target docker.service qwen-benchmark-api.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ecs-user
Group=ecs-user
SupplementaryGroups=docker
WorkingDirectory=/srv/qwen-benchmark
EnvironmentFile=/srv/qwen-benchmark/config/benchmark.env
EnvironmentFile=/srv/qwen-benchmark/config/benchmark.secret.env
ExecStart=/srv/qwen-benchmark/venv/bin/qwen-benchmark-worker
Restart=on-failure
RestartSec=5
TimeoutStopSec=90
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=tmpfs
CacheDirectory=qwen-benchmark/harbor
BindPaths=/var/cache/qwen-benchmark/harbor:/home/ecs-user/.cache/harbor
ReadWritePaths=/srv/qwen-benchmark /var/lib/docker /var/run/docker.sock

[Install]
WantedBy=multi-user.target
