# benchmark-service/deploy/nginx-supervisor.conf
[program:qwen-benchmark-nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/qwen-benchmark/logs/nginx-supervisor.log
stderr_logfile=/tmp/qwen-benchmark/logs/nginx-supervisor-error.log
priority=40
