# benchmark-service/deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/qwen-benchmark/logs/supervisord.log
pidfile=/tmp/qwen-benchmark/supervisord.pid

[program:dockerd]
command=/usr/bin/dockerd --host=unix:///var/run/docker.sock --data-root=/tmp/qwen-benchmark/docker --storage-driver=overlay2
autorestart=true
startsecs=5
stdout_logfile=/tmp/qwen-benchmark/logs/dockerd-supervisor.log
stderr_logfile=/tmp/qwen-benchmark/logs/dockerd-supervisor-error.log
priority=10

[program:benchmark-api]
command=/bin/bash -lc 'set -a; source /mnt/workspace/qwen-benchmark/config/benchmark.env; exec /mnt/workspace/qwen-benchmark/venv/bin/qwen-benchmark-api'
directory=/mnt/workspace/qwen-benchmark
autorestart=true
startsecs=3
stdout_logfile=/tmp/qwen-benchmark/logs/api.log
stderr_logfile=/tmp/qwen-benchmark/logs/api-error.log
priority=20

[program:benchmark-worker]
command=/bin/bash -lc 'set -a; source /mnt/workspace/qwen-benchmark/config/benchmark.env; exec /mnt/workspace/qwen-benchmark/venv/bin/qwen-benchmark-worker'
directory=/mnt/workspace/qwen-benchmark
autorestart=true
startsecs=3
stdout_logfile=/tmp/qwen-benchmark/logs/worker.log
stderr_logfile=/tmp/qwen-benchmark/logs/worker-error.log
priority=30
